# nginx/crontab
0 12 * * * certbot renew --quiet && nginx -s reload
